# infra/trading-dashboard.service
[Unit]
Description=Trading Dashboard Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/damiverse_apps/trading_app
ExecStart=/root/damiverse_apps/trading_app/.venv/bin/python -m apps.dashboard.serve_dashboard
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=TD_HOST=127.0.0.1
Environment=TD_PORT=8050
Environment=TRADING_APP_ROOT=/root/damiverse_apps/trading_app

[Install]
WantedBy=multi-user.target
